# raspmidi-lcd.service
[Unit]
Description=RaspMIDI LCD Service
After=network.target raspmidi.service
Wants=raspmidi.service

[Service]
Type=simple
User=matheus
Group=matheus
WorkingDirectory=/home/matheus/RaspMIDI
Environment=DISPLAY=:0
Environment=PYTHONPATH=/home/matheus/RaspMIDI
ExecStart=/home/matheus/RaspMIDI/venv/bin/python /home/matheus/RaspMIDI/app/lcd_service.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target 